# Buildbot-worker/conf/supervisor/worker.conf
[program:worker]
command=/root/venv/buildbot-worker start /root/worker
autostart=true
autorestart=true
stdout_logfile = /var/log/buildbot.log
