# QorIQ/conf/supervisor/rpcbind.conf
[program:rpcbind]
command=/sbin/rpcbind -i -w
autostart=true
autorestart=true
