# Letschat/conf/supervisor/letschat.ini
[program:letschat]
directory = /root/lets-chat
command = node app.js
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/letschat.log
