# Kallithea/conf/supervisor/kallithea.ini
[program:kallithea]
directory = /root
command = /root/venv/bin/paster serve /var/lib/kallithea/data/production.ini
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/kallithea.log
