# Janus/conf/supervisor/janus.conf
[program:janus-gateway]
command = janus
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/janus-gateway.log
