# Trac/conf/supervisor/trac.ini
[program:trac]
directory = /root
command = /root/venv/bin/python /root/bootstrap.py trac-plugin-test 8001
#autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/trac.log
