# Jenkins-master/conf/supervisor/jenkins.ini
[program:jenkins-master]
command=jenkins.sh
user=jenkins
environment=JENKINS_OPTS="--prefix=/jenkins"
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/jenkins.log
stopsignal=QUIT
