# Buildbot-master/conf/supervisor/master.ini
[program:master]
command = /root/venv/bin/buildbot start /var/lib/buildbot/master
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/buildbot.log
