# Errbot/conf/supervisor/errbot.ini
[program:errbot]
directory = /root/errbot-moogle
command = /root/venv/bin/errbot
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/errbot.log
